# root/etc/services.d/php-fpm/run
#!/command/with-contenv bash

exec /usr/sbin/php-fpm82 -F
