# root/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
fdmove -c 2 1

foreground { echo "testing nginx config" }
foreground { /usr/sbin/nginx -t -c /config/nginx/nginx.conf }

foreground { echo "starting nginx" }
exec /usr/sbin/nginx -c /config/nginx/nginx.conf

